# setup/configs/mywebapp.socket.template
[Unit]
Description=My Web App socket

[Socket]
ListenStream=${APP_PORT}
Accept=no

[Install]
WantedBy=sockets.target
